# Session_29/Task_1/ansible/roles/flask_app/templates/service_name.service.j2
[Unit]
Description={{ service_description }}
Wants=network-online.target
After=network-online.target

[Service]
WorkingDirectory={{ working_dir }}
ExecStart={{ exec_start }}
Restart=always
CPUQuota={{ service_cpu }}
MemoryMax={{ service_memory }}

[Install]
WantedBy=multi-user.target
